# battery-alarm.target
[Unit]
Description=Battery Alarm Target
StopWhenUnneeded=yes
